# .crontab
1,31 * * * * ~/bin/download-bing-wallpaper.js
